# functions/officeToPdf/supervisord.conf
; functions/officeToPdf/supervisord.conf
;
; (2026-05-22 EOD+2 — T4) Start two co-located processes inside the same
; Cloud Run container:
;   - gotenberg → localhost:3000 (private; only the Functions handler talks to it)
;   - functions  → $PORT=8080 (Cloud Run public — receives Storage events via Eventarc)
;
; Both processes run as gotenberg user (non-root) per the Gotenberg base image's
; security posture.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:gotenberg]
command=/usr/bin/gotenberg --api-port=3000
user=gotenberg
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:functions]
command=/usr/bin/npx @google-cloud/functions-framework --target=officeToPdf --signature-type=cloudevent --port=8080
directory=/workspace
user=gotenberg
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
